# rt-cron/crontab
* * * * * fetchmail -f /home/rt4/fetchmailrc
